# PrintFlowLite/packages/docker/supervisord.conf
# =============================================================================
# Supervisor configuration for PrintFlowLite Docker container
# =============================================================================
# Manages multiple processes: CUPS, PrintFlowLite server, and XML-RPC notifier.
# =============================================================================

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

# -----------------------------------------------------------------------------
# CUPS Print Server
# -----------------------------------------------------------------------------
[program:cupsd]
command=/usr/sbin/cupsd -f
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/cupsd-stdout.log
stderr_logfile=/var/log/supervisor/cupsd-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
user=root
priority=100

# Enable remote CUPS administration
[program:cupsctl]
command=/bin/bash -c "cupsctl --remote-any 2>/dev/null || true"
autostart=true
autorestart=false
stdout_logfile=/var/log/supervisor/cupsctl.log
stderr_logfile=/var/log/supervisor/cupsctl.log
user=root
priority=200

# -----------------------------------------------------------------------------
# PrintFlowLite CUPS Notifier (runs as root for XML-RPC)
# -----------------------------------------------------------------------------
[program:printflowlite-notifier]
command=/opt/printflowlite/providers/cups/linux-x64/roottasks
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/notifier-stdout.log
stderr_logfile=/var/log/supervisor/notifier-stderr.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
user=root
priority=300

# -----------------------------------------------------------------------------
# PrintFlowLite Server
# -----------------------------------------------------------------------------
[program:printflowlite]
command=/opt/printflowlite/server/bin/linux-x64/app-server start
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/printflowlite-stdout.log
stderr_logfile=/var/log/supervisor/printflowlite-stderr.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
user=printflowlite
priority=500
# Wait for CUPS to be ready first
startsecs=30
startretries=3

# -----------------------------------------------------------------------------
# Log rotation for supervisor logs
# -----------------------------------------------------------------------------
[program:logrotate]
command=/usr/sbin/logrotate /etc/logrotate.d/supervisor
autostart=false
autorestart=false
stdout_logfile=/var/log/supervisor/logrotate.log
stderr_logfile=/var/log/supervisor/logrotate.log
